# job-board/django-backend/docker/supervisor/supervisord.conf
# docker/supervisor/supervisord.conf - Process management for production

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# Gunicorn Django Application
[program:gunicorn]
command=gunicorn config.wsgi:application --bind 0.0.0.0:8000 --workers 4 --worker-class gevent --worker-connections 1000 --max-requests 1000 --max-requests-jitter 100 --timeout 30 --keep-alive 5 --preload
directory=/app
user=django
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/gunicorn.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/gunicorn_error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=5

# Daphne for WebSocket support (optional)
[program:daphne]
command=daphne -b 0.0.0.0 -p 8001 config.asgi:application
directory=/app
user=django
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/daphne.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/daphne_error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=5

# Nginx Web Server
[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/nginx.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/nginx_error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=5

# Log rotation
[program:logrotate]
command=/usr/sbin/logrotate /etc/logrotate.conf
autostart=false
autorestart=false
startsecs=0